# cfan.service
[Unit]
Description=cfan, smart fan controlling software

[Service]
Type=simple
ExecStart=/usr/local/bin/cfan

[Install]
WantedBy=multi-user.target